# supervisord.conf
; ============================================================================
; supervisord — process orchestration for the production Railway container.
; ============================================================================
; Runs php-fpm + nginx side by side in the same container. This file is
; copied into /etc/supervisor/supervisord.conf during the Docker build and
; invoked as PID 1 in the container CMD.
;
; Why this exists: production used to launch `php artisan serve`, PHP's
; built-in single-threaded dev server. Replacing it with a real php-fpm
; worker pool fronted by nginx is the single biggest perf win for the app
; (measured: under just 6 concurrent requests, TTFB grew from 213 ms on
; the first to 795 ms on the fifth with the old setup — that's gone now).
; ============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/run/supervisord.pid

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=5

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10
; ============================================================================
; queue-worker — Tier-2 async pipeline
; ----------------------------------------------------------------------------
; Drains the `high` and `default` queues from the database driver (jobs
; table). Single concurrency by design (1 replica × 1 worker process).
; Per-send pacing is NOT provided by single concurrency alone — a worker
; with no jobs queued still fires sends back-to-back at network speed — so
; the image job applies an explicit WHATSAPP_SEND_PACING_MS sleep between
; Meta calls (config/whatsapp.php) to replace the old in-request sleep(1).
;
;  --tries=5 / --backoff=15,60,180,300 : 5 attempts on a long backoff; the
;                                    job also long-releases on Meta throttle
;                                    codes so retries land after the window.
;  --timeout=120                   : hard per-job cap; kept below the
;                                    connection retry_after (180) so a reaped
;                                    attempt is retried, not double-run.
;  --max-time=3600                 : worker self-restarts every hour so a
;                                    Railway deploy can't strand it with
;                                    stale bytecode.
;  --queue=high,default            : `high` drained first (admin approve,
;                                    admin/customer resend, webhook
;                                    reply). `default` reserved for
;                                    lower-priority work that doesn't
;                                    exist yet.
;  --sleep=3                       : poll the jobs table every 3 s when
;                                    idle. Default value; trade-off
;                                    between first-job latency and DB
;                                    poll rate.
;
; IMPORTANT: production must run QUEUE_CONNECTION=database so dispatch()
; uses this worker. With QUEUE_CONNECTION=sync every dispatch() runs inline
; on the webhook request thread, so a single tap on a large booking tries
; to send every ticket image within one HTTP request and is truncated when
; the request times out — the remaining tickets stay delivery_status=
; 'pending' and never send until another tap. Only fall back to `sync` as
; a deliberate emergency rollback, never as the steady state.
; ============================================================================
[program:queue-worker]
command=php /app/artisan queue:work --tries=5 --backoff=15,60,180,300 --timeout=120 --max-time=3600 --sleep=3 --queue=high,default
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
