# supervisord.conf
; ============================================================================
; supervisord — process orchestration for the production Railway container.
; ============================================================================
; Runs php-fpm + nginx side by side in the same container. This file is
; copied into /etc/supervisor/supervisord.conf during the Docker build and
; invoked as PID 1 in the container CMD.
;
; Why this exists: production used to launch `php artisan serve`, PHP's
; built-in single-threaded dev server. Replacing it with a real php-fpm
; worker pool fronted by nginx is the single biggest perf win for the app
; (measured: under just 6 concurrent requests, TTFB grew from 213 ms on
; the first to 795 ms on the fifth with the old setup — that's gone now).
; ============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/run/supervisord.pid

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=5

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10
